# portfolio.service
[Unit]
Description=Portfolio Application
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=forking
User=www-data
Group=www-data
WorkingDirectory=/var/www/portfolio
ExecStart=/usr/local/bin/docker-compose up -d
ExecStop=/usr/local/bin/docker-compose down
ExecReload=/usr/local/bin/docker-compose restart
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target